# deploy/smb-watch.timer
[Unit]
Description=smb-watch を定期実行する (平日 9〜17 時 JST、毎正時に 1 回)

[Timer]
# 平日 (月〜金) の 9:00〜17:00 JST を毎正時に実行 = 1 日 9 回 (JST 9,10,…,17 時)。
#
# OnCalendar はホスト TZ で解釈される。ohishi-data は **host TZ = UTC** 運用で、
# 他に時刻依存の cron (backup / update_ichi) があるため host TZ は変更しない。
# JST 9〜17 時 = UTC 0〜8 時 (曜日も Mon-Fri で一致) なので、ここは **UTC で記述**する。
# = JST 9:00→UTC 0:00 / JST 17:00→UTC 8:00。host TZ を JST にした場合はこの値を
#   09..17 に直すこと (systemd 245 は OnCalendar の TZ 接尾辞 v252+ が使えない)。
# マシン停止中に逃した分は起動後にまとめて 1 回走る (Persistent)。
OnCalendar=Mon-Fri *-*-* 00..08:00:00
Persistent=true
# 複数ホストでの同時刻集中を避けるためのランダム遅延 (最大 5 分)。
RandomizedDelaySec=300
Unit=smb-watch.service

[Install]
WantedBy=timers.target
